# reins/scraper/com.reins.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.reins.scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/masato/test-project/reins/scraper/venv/bin/python</string>
        <string>/Users/masato/test-project/reins/scraper/reins_scraper.py</string>
    </array>

    <!-- 毎日8時・12時・16時に実行 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/masato/test-project/reins/scraper/reins_scraper.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/masato/test-project/reins/scraper/reins_scraper.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/masato/test-project/reins/scraper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HEADLESS</key>
        <string>false</string>
    </dict>
</dict>
</plist>
